# sample/upstart.conf
description "hydra - httpd multiplexer"
author "Richard Whitehouse <richard@whiuk.com>"

start on started networking
stop on runlevel [!2345]

exec /usr/local/sbin/hydra

